# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:mailhog]
command=/usr/local/bin/MailHog

[program:hogbridge]
command=/usr/local/bin/HogBridge
